# archlinux/wolframe.install
post_install() {
  # Don't enable Wolframe server at install time, just inform root how this is done
echo
echo Use 'systemctl enable wolframed' to enable the server at startup
echo
}
